# roles/proces-exporter/files/process-exporter.service
[Unit]
Description=Process Exporter for Prometheus
Documentation=https://github.com/ncabatoff/process-exporter
Wants=network-online.target
After=network-online.target

[Service]
User=process-exporter
Group=process-exporter
Type=simple
Restart=on-failure
ExecStart=/usr/local/bin/process-exporter-0.6.0.linux-amd64/process-exporter \
  --config.path /etc/process-exporter/process-exporter.yml \
  --web.listen-address=:9256

[Install]
WantedBy=multi-user.target